# roman-os-system-files/etc/systemd/system/roman-os-pci-latency.service
[Unit]
Description=Optimize PCI Latency Timers for Responsiveness
Documentation=https://wiki.archlinux.org/title/PCI_latency
After=multi-user.target

[Service]
Type=oneshot
User=root
ExecStart=/usr/local/bin/kiro-pci-latency
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kiro-pci-latency

# Continue boot even if script fails (kiro-pci-latency changes are non-critical)
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
