# roman-os-system-files/etc/udev/rules.d/62-network-optimization.rules
# ============================================================================
# Network Interface Optimization - Enhance responsiveness
# ============================================================================
# Optimizes network drivers for low latency and responsiveness

# Intel server NICs: disable interrupt coalescing for low latency
# Prevents buffering delays in network packets
# ~5-10ms latency reduction for desktop/gaming
# add|move: the interface is renamed (ethN -> enpXsY) via a move uevent after
# add; matching move runs ethtool against the final name. The `|| true` guard
# swallows the stale-name race on add so udev logs no "Process failed" warning.
ACTION=="add|move", SUBSYSTEM=="net", DRIVERS=="igb|ixgbe|i40e", \
  RUN+="/bin/sh -c '/usr/bin/ethtool -C $name rx-usecs 0 rx-frames 0 tx-usecs 0 tx-frames 0 || true'"

# Intel client NICs (e1000e: I217/I218/I219): only rx-usecs is supported;
# rx-frames/tx-usecs/tx-frames return EINVAL (ethtool exit 1) on these controllers
ACTION=="add|move", SUBSYSTEM=="net", DRIVERS=="e1000e", \
  RUN+="/bin/sh -c '/usr/bin/ethtool -C $name rx-usecs 0 || true'"

# Realtek Ethernet: balanced power and performance
# Uses minimal coalescing (64us) to reduce latency while saving power
ACTION=="add", SUBSYSTEM=="net", DRIVERS=="r8169", \
  RUN+="/usr/bin/ethtool -C $name rx-usecs 64 rx-frames 1"

# Broadcom: Similar to Realtek - balanced approach
ACTION=="add", SUBSYSTEM=="net", DRIVERS=="bcm63xx_enet|tg3|bnx2x", \
  RUN+="/usr/bin/ethtool -C $name rx-usecs 32 rx-frames 1"

# Mellanox/ConnectX: No coalescing for ultra-low latency
ACTION=="add", SUBSYSTEM=="net", DRIVERS=="mlx5_core|mlx4_en", \
  RUN+="/usr/bin/ethtool -C $name rx-usecs 0 tx-usecs 0"

# WiFi: enable 802.11n/802.11ac for better throughput
ACTION=="add", SUBSYSTEM=="net", KERNEL=="w*", \
  RUN+="/usr/bin/iw $name set txpower auto"

# Reduce MTU jumbogram issues (common on WiFi)
# Standard MTU 1500 is safest for wireless networks
ACTION=="add", SUBSYSTEM=="net", KERNEL=="wlan*|wif*", ATTR{mtu}="1500"

# Enable Generic Segmentation Offload (GSO) for high-speed interfaces
# Helps with TCP throughput on gigabit+ connections
# (e1000e omitted: GSO is already on by default there; forcing it returns
#  EOPNOTSUPP — ethtool exit 92 — on consumer I219-class controllers)
ACTION=="add|move", SUBSYSTEM=="net", DRIVERS=="ixgbe|i40e", \
  RUN+="/bin/sh -c '/usr/bin/ethtool -K $name gso on || true'"

# Disable TSO if causing latency issues (uncomment if needed)
# Most modern drivers handle TSO well, but some problematic ones may need this
# ACTION=="add", SUBSYSTEM=="net", DRIVERS=="your_driver", \
#   RUN+="/usr/bin/ethtool -K $name tso off gso off"
