# roman-os-system-files/etc/udev/rules.d/67-laptop-optimization.rules
# ============================================================================
# Laptop/Battery Optimization - Power efficiency
# ============================================================================
# Critical for laptop users running keyboard-driven, minimal WMs

# Disabled: Not applicable to desktop systems
# # Disable bluetooth by default (power save)
# ACTION=="add", SUBSYSTEM=="pci", ATTRS{vendor}=="0x8087", ATTR{power/control}="auto"
#
# # Reduce WiFi power consumption when idle
# ACTION=="add", SUBSYSTEM=="net", KERNEL=="wlan*", \
#   RUN+="/sbin/iw $name set power_save on"
#
# # Lower CPU frequency scaling for battery
# ACTION=="add", SUBSYSTEM=="dmi", ATTR{id}=="", \
#   RUN+="/bin/bash -c 'echo powersave > /sys/devices/system/cpu/cpu*/cpufreq/scaling_governor 2>/dev/null || true'"
